# debian/postinst
systemctl daemon-reload
systemctl enable ukui-touch.service
systemctl restart ukui-touch.service
