# debian/prerm
systemctl stop ukui-touch.service
systemctl disable ukui-touch.service
